# deploy/coldp-editor.service
# systemd unit for the ColDP editor backend (Spring Boot jar) on the GBIF dev VM.
# Install to /etc/systemd/system/coldp-editor.service, then:
#   sudo systemctl daemon-reload && sudo systemctl enable --now coldp-editor
# Secrets (DB password, ORCID) come from the EnvironmentFile — never hard-code them here.
[Unit]
Description=ColDP Editor backend (Spring Boot, JDK 25)
After=network-online.target postgresql.service
Wants=network-online.target

[Service]
User=col
Group=col
EnvironmentFile=/etc/coldp-editor/coldp-editor.env
# JAVA_OPTS comes from the env file; the jar is deployed by deploy.sh (see README).
ExecStart=/usr/bin/env java $JAVA_OPTS -jar /opt/coldp-editor/coldp-editor-backend.jar
SuccessExitStatus=143
Restart=on-failure
RestartSec=5
# Spring Boot listens on 127.0.0.1:8080 (Apache proxies to it); it is not exposed directly.

[Install]
WantedBy=multi-user.target
